# deploy/macos/com.mino.rebalancing.status.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mino.rebalancing.status</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/myno-bot/Desktop/bot &amp;&amp; exec /Users/myno-bot/Desktop/bot/.venv/bin/python -m rebalancing.status_server</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/myno-bot/Desktop/bot</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>/Users/myno-bot/Desktop/bot/src</string>
    <key>DATABASE_URL</key>
    <string>postgresql://rebalancing@127.0.0.1:5432/rebalancing_learning</string>
    <key>ENGINE_HOST</key>
    <string>127.0.0.1</string>
    <key>ENGINE_PORT</key>
    <string>8788</string>
    <key>ENGINE_CORS_ORIGIN</key>
    <string>*</string>
    <key>ENGINE_WEBHOOK_TOKEN_FILE</key>
    <string>/Users/myno-bot/Desktop/bot/.state/engine_webhook_token</string>
    <key>TV_WEBHOOK_PASSPHRASE_FILE</key>
    <string>/Users/myno-bot/Desktop/bot/.state/tradingview_webhook_passphrase</string>
    <key>ENGINE_SIGNAL_STORE_PATH</key>
    <string>/Users/myno-bot/Desktop/bot/.state/tradingview_alerts.json</string>
    <key>ENGINE_SIGNAL_STORE_MAX_RECORDS</key>
    <string>200</string>
    <key>ENGINE_STATUS_ALERT_EVENT_LIMIT</key>
    <string>200</string>
    <key>ENGINE_TV_MAX_ALERT_AGE_SECONDS</key>
    <string>0</string>
    <key>ENGINE_TV_ALLOW_PROBE_ENTRIES</key>
    <string>true</string>
    <key>ENGINE_TV_PROBE_LEVERAGE</key>
    <string>0.5</string>
    <key>PAPER_TRADING_ENABLED</key>
    <string>true</string>
    <key>PAPER_INITIAL_EQUITY</key>
    <string>1000</string>
    <key>PAPER_STATE_PATH</key>
    <string>/Users/myno-bot/Desktop/bot/.state/paper_trading.json</string>
    <key>PAPER_MAX_LEVERAGE</key>
    <string>2</string>
    <key>PAPER_MIN_ORDER_NOTIONAL</key>
    <string>10</string>
    <key>PAPER_MIN_REBALANCE_MINUTES</key>
    <string>60</string>
    <key>PAPER_BLOCK_NEW_ENTRIES_LOSS_PCT</key>
    <string>-1.0</string>
    <key>PAPER_FEE_RATE</key>
    <string>0.0004</string>
    <key>PAPER_SLIPPAGE_BPS</key>
    <string>0</string>
    <key>PAPER_TRADE_HISTORY_LIMIT</key>
    <string>500</string>
    <key>FALLBACK_EQUITY</key>
    <string>1000</string>
    <key>FALLBACK_WALLET_BALANCE</key>
    <string>1000</string>
    <key>DAY_START_EQUITY</key>
    <string>1000</string>
    <key>WEEK_START_EQUITY</key>
    <string>1000</string>
    <key>MONTH_START_EQUITY</key>
    <string>1000</string>
    <key>BINANCE_LIVE_TRADING</key>
    <string>false</string>
    <key>BINANCE_ENABLE_ORDERS</key>
    <string>false</string>
    <key>BINANCE_POSITION_MODE</key>
    <string>ONE_WAY</string>
    <key>MARKET_INTERNALS_CACHE_SECONDS</key>
    <string>300</string>
    <key>MARKET_INTERNALS_UNIVERSE_LIMIT</key>
    <string>200</string>
    <key>MARKET_INTERNALS_BREADTH_LIMIT</key>
    <string>100</string>
    <key>MARKET_INTERNALS_BREADTH_WORKERS</key>
    <string>8</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/myno-bot/Desktop/bot/.state/status_server.launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/myno-bot/Desktop/bot/.state/status_server.launchd.err.log</string>
</dict>
</plist>
